# saltstack/salt/files/openhab/openhab-minute.timer
[Unit]
Description=Openhab in-minute script timer

[Timer]
OnCalendar=*:*
Unit=openhab-minute.service

[Install]
WantedBy=timers.target
