# saltstack/salt/files/upgrade-checker/services-upgrade-checker.timer
[Unit]
Description=Updates-checker script timer

[Timer]
OnCalendar=weekly
Unit=services-upgrade-checker.service

[Install]
WantedBy=timers.target
