# saltstack/salt/files/home/openhab-hourly.timer
[Unit]
Description=Openhab hourly script timer

[Timer]
OnCalendar=hourly
Unit=openhab-hourly.service

[Install]
WantedBy=timers.target
