# saltstack/salt/files/pws-pve/pve-backup.service
[Unit]
Description=Backup PVE service

[Service]
Type=notify
RemainAfterExit=no
ExecStart=/opt/backup/bin/backup
TimeoutStartSec=0
